# discord-canary-installer.timer
[Unit]
Description=Updating Discord Canary

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target