# discord-installer.timer
[Unit]
Description=Updating Discord

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target