# discord-canary-installer.service
[Unit]
Description=Updating Discord Canary
Requires=network-online.target
After=network-online.target

[Service]
Type=oneshot
RemainAfterExit=no
ExecStart=/usr/libexec/discord-installer/build.sh DiscordCanary discord-canary https://discordapp.com/api/download/canary?platform=linux&format=tar.gz

[Install]
WantedBy=graphical.target